# example-xremap.service
[Unit]
Description=xremap user service
After=default.target

[Service]
ExecStart=/usr/bin/xremap %h/.config/xremap/config.yml
Restart=always
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
