# doc/kde-wayland/kde-socket-helper.service
[Unit]
Description=Xremap KDE Socket Helper
After=graphical-session.target
Requires=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/bin/kde-socket-helper
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical-session.target
